# rootfs/etc/s6-overlay/s6-rc.d/wg0/up
#!/usr/bin/execlineb -P

# start the vpn
wg-quick up /config/wg0.conf
