# rootfs/etc/s6-overlay/s6-rc.d/logs/run
#!/usr/bin/execlineb -P

# listen for logs
tail -F -n0 /config/qBittorrent/data/logs/qbittorrent.log
